# config/etc/udev/rules.d/backlight.rules
RUN+="/bin/chgrp video /sys/class/backlight/intel_backlight/brightness"
RUN+="/bin/chmod g+w /sys/class/backlight/intel_backlight/brightness"
RUN+="/bin/chown 1000 /sys/class/backlight/intel_backlight/brightness"
